# config/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/data/logs/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:mihomo]
command=/usr/local/bin/mihomo -d /data/mihomo
autorestart=true
startsecs=3
priority=10
stdout_logfile=/data/logs/mihomo-stdout.log
stderr_logfile=/data/logs/mihomo-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

[program:wg-easy]
command=/scripts/start-wg-easy.sh
autorestart=true
startsecs=5
priority=20
stdout_logfile=/data/logs/wg-easy-stdout.log
stderr_logfile=/data/logs/wg-easy-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

# Long-running daemon — re-applies iptables FORWARD rules for wg0
# if they get flushed (e.g. Docker restart). Mihomo handles routing.
[program:setup-routing]
command=/scripts/setup-routing.sh
autorestart=true
startsecs=0
priority=30
stdout_logfile=/data/logs/routing-stdout.log
stderr_logfile=/data/logs/routing-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
